# deploy/blackstate.service
[Unit]
Description=BlackState Casino (Bun)
After=network.target

[Service]
Type=simple
User=blackstate
WorkingDirectory=/opt/blackstate
EnvironmentFile=/opt/blackstate/.env
# Vérifie le chemin avec : sudo -u blackstate which bun
ExecStart=/home/blackstate/.bun/bin/bun run server.ts
Restart=always
RestartSec=3

# Durcissement basique
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
